# upgrades/upgrade3/+/etc/rc.local
#!/bin/sh
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# update mode if needed
MODE=`/sbin/checkswitch`
[ "$MODE" != 'none' ] && /sbin/modechange $MODE

# setup ap access check
/sbin/fix_sta.sh > /dev/null 2>&1 &

# start internet check to light LED
/sbin/conn_check &

exit 0
