# root/etc/rc.local
#!/bin/sh
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/sbin/fix_sta.sh > /dev/null 2>&1 &

exit 0
